# com.openclaw.linkedin-writer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.openclaw.linkedin-writer</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/node</string>
		<string>/Users/timothyryan/.openclaw/workspace/scripts/generate-linkedin-post.js</string>
	</array>
	<key>StartCalendarInterval</key>
	<array>
		<!-- Tuesday 8 AM EST -->
		<dict>
			<key>Hour</key>
			<integer>8</integer>
			<key>Minute</key>
			<integer>0</integer>
			<key>Weekday</key>
			<integer>3</integer>
		</dict>
		<!-- Wednesday 8 AM EST -->
		<dict>
			<key>Hour</key>
			<integer>8</integer>
			<key>Minute</key>
			<integer>0</integer>
			<key>Weekday</key>
			<integer>4</integer>
		</dict>
		<!-- Thursday 8 AM EST -->
		<dict>
			<key>Hour</key>
			<integer>8</integer>
			<key>Minute</key>
			<integer>0</integer>
			<key>Weekday</key>
			<integer>5</integer>
		</dict>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>BRAVE_API_KEY</key>
		<string>BSAHJ3Wmk1IbHNqEsACADrcFLfW5eLc</string>
		<key>HOME</key>
		<string>/Users/timothyryan</string>
	</dict>
	<key>StandardOutPath</key>
	<string>/Users/timothyryan/.openclaw/workspace/.linkedin-writer.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/timothyryan/.openclaw/workspace/.linkedin-writer-error.log</string>
	<key>RunAtLoad</key>
	<false/>
	<key>TimeZone</key>
	<string>America/New_York</string>
</dict>
</plist>
