# jtsa-backend/scripts/com.jtsa.opus-reviewer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jtsa.opus-reviewer</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/timothyryan/.openclaw/workspace/jtsa-backend/scripts/trigger-review.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>StandardOutPath</key>
    <string>/Users/timothyryan/.openclaw/workspace/jtsa-backend/logs/reviewer-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/timothyryan/.openclaw/workspace/jtsa-backend/logs/reviewer-launchd.err</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
